# telegram-claude.service
[Unit]
Description=Telegram Claude Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/user/telegram-claude
EnvironmentFile=/home/user/telegram-claude/.env
ExecStart=/home/user/.bun/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5
Environment=PATH=/home/user/.local/bin:/home/user/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
